# hipchat/ngrok.sv.conf
[program:ngrok]
command=/usr/bin/ngrok -log=stdout 3000
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

