# hipchat/ngrok.sv.conf
[program:ngrok]
command=/usr/bin/ngrok -log=stdout 3000
redirect_stderr=true
stdout_logfile=/var/log/supervisor/ngrok.log

